# bpi/lib/systemd/system/bthelper@.service
[Unit]
Description=Bluetooth hci up helper
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bthelper %I
